# router/scripts/com.jarvis.router.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Example LaunchAgent for the Jarvis router.

  Install (setup.sh does this automatically on macOS):
    sed "s|__HOME__|$HOME|g; s|__NODE_BIN__|$(command -v node)|g; s|__NODE_DIR__|$(dirname $(command -v node))|g" \
      com.jarvis.router.plist.example \
      > ~/Library/LaunchAgents/com.jarvis.router.plist
    launchctl load ~/Library/LaunchAgents/com.jarvis.router.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jarvis.router</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__HOME__/.claude/jarvis/router/node_modules/.bin/tsx</string>
    <string>src/index.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__HOME__/.claude/jarvis/router</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/.claude/jarvis/logs/router.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.claude/jarvis/logs/router.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__NODE_DIR__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
